# deploy/playbooks/templates/[[app_name]].service.j2
[Unit]
Description={{ app_name }} (Bun / Elysia + SPA)
After=network.target

[Service]
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ server_current }}
EnvironmentFile=/etc/{{ app_name }}.env
Environment=NODE_ENV=production
ExecStart=/root/.bun/bin/bun ./index.ts
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
